# rclone/rootfs/etc/services.d/unionfs-mount/finish
#!/usr/bin/with-contenv sh

#set -x

# make sure unmounted
/usr/bin/fusermount -u -z ${MountPoint}
/usr/bin/umount ${MountPoint}
