# rclone/rootfs/etc/services.d/unionfs-mount/finish
#!/usr/bin/with-contenv sh

#set -x

# make sure unmounted
fusermount -u -z ${MountPoint}
umount ${MountPoint}
