# filebot/rootfs/etc/services.d/filebot-watcher/run
#!/usr/bin/with-contenv bash

# run script as user abc
exec s6-setuidgid abc \
  /bin/bash /app/inotify-watch.sh
